# templates/etc/supervisor/conf.d/notes.conf.erb
[program:notes_script]
directory=/var/www/notes
command=/var/www/notes/bin/python manage.py rqworker
stopasgroup=true
stopsignal=QUIT
stderr_logfile=/var/log/notes.err.log
stdout_logfile=/var/log/notes.out.log

